# ansible/roles/git/local_config/templates/gpg-preload-passphrase.service.j2
[Unit]
Description=Preload GPG passphrase into gpg-agent from GNOME Keyring
# graphical-session.target is reached after the desktop (and keyring) are up.
# Using PartOf means this unit stops/restarts whenever the graphical session does.
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/gpg-preload-passphrase.sh
RemainAfterExit=no

[Install]
WantedBy=graphical-session.target
